# webapp/flask-backend/magic-webapp-dev.conf
[supervisord]
logfile=/home/ubuntu/magic-webapp/flask-backend/logs/supervisor.log
pidfile=/home/ubuntu/magic-webapp/flask-backend/logs/supervisor.pid
nodaemon=true ; run supervisord as a daemon when the system starts

[program:magic-webapp]
command=/home/ubuntu/.virtualenvs/magic-flask/bin/gunicorn -c gunicorn_config.py app:app
directory=/home/ubuntu/magic-webapp/flask-backend
environment=ENV_FILE=development.py
user=ubuntu
redirect_stderr=true
autostart=true
autorestart=true
stopasgroup=true  ; send the stop signal to the whole process group
stopsignal=QUIT   ; change to a signal that your app handles correctly
killasgroup=true  ; kill the whole process group if the app doesn't stop in time
stopwaitsecs=10   ; how long to wait before forcibly stopping

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface